# cron.d/crontab
* * * * *	www-data	bash -c "/usr/bin/php -r 'echo \"Hello\";'  > /dev/stdout 2>/dev/stderr"