# cron.d/crontab
* * * * *   bash -c "/usr/bin/php -r 'echo \"Hello\";'"