# cron.d/crontab
* * * * *   /usr/bin/php -r 'echo getenv("PHP_VERSION");'
* * * * *   /usr/bin/php -r 'echo "Hello";'